# scripts/e2e/ssh-target/podup-e2e-noop.service
[Unit]
Description=pod-upgrade-trigger E2E noop (no side effects)

[Service]
Type=simple
ExecStart=/usr/bin/sleep infinity

[Install]
WantedBy=default.target
